# assistant/supervisord.conf
[program:rasa]
command=rasa run --cors "*" --model models/20210907-184826.tar.gz

[program:rasa-actions]
command=rasa run actions --cors "*"